# setup_servidor/com.arkaitz.bot_datos.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.arkaitz.bot_datos</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/mac/Desktop/Arkaitz/telegram_bot_datos/bot_datos.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/mac/Desktop/Arkaitz/telegram_bot_datos</string>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/Users/mac/Desktop/Arkaitz/logs/bot_datos.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mac/Desktop/Arkaitz/logs/bot_datos.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
